# systemd/cloudflared-template.service
[Unit]
Description=Cloudflare Tunnel - {{TUNNEL_NAME}}
After=network.target

[Service]
Type=notify
User=root
# Fix ping_group_range to enable ICMP proxy (prevents GID warnings)
ExecStartPre=/bin/bash -c 'echo "0 2000" > /proc/sys/net/ipv4/ping_group_range || true'
ExecStart=/usr/local/bin/cloudflared tunnel --config /etc/cloudflared/{{TUNNEL_NAME}}.yml run
Restart=on-failure
RestartSec=5s
KillMode=mixed
KillSignal=SIGINT
TimeoutStopSec=30s

[Install]
WantedBy=multi-user.target
